# accel-sensor-bridge.service
[Unit]
Description=IIO Accelerometer/IMU Sensor Bridge
After=network.target

[Service]
ExecStart=/usr/local/bin/accel_sensor_bridge.py
Restart=on-failure
RestartSec=30s
StartLimitIntervalSec=60
StartLimitBurst=2

User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target